# files/upstart/chef-server.conf
start on runlevel [2345]
stop on runlevel [06]

respawn

script
  su - USER -c "chef-server -N -e production"
end script
