# files/upstart/chef-server-webui.conf
start on started chef-server
stop on stopped chef-server

respawn

script
  su - USER -c "rvmsudo chef-server-webui -p 4040 -e production"
end script
