# files/upstart/chef-solr.conf
start on runlevel [2345]
stop on runlevel [06]

respawn

script
  su - USER -c "rvmsudo chef-solr"
end script
